# crontab
*/1 * * * * python /app/app/daemon.py | tee -a /app/logfile >/dev/null
